# check_active_session.service
[Unit]
Description = Check defis active session

[Service]
ExecStart=/etc/.defis_adh_X11/.defis_active_session

[Install]
WantedBy=default.target

